# packaging/deb/postinst
#!/bin/sh
# postinst — install the default config on first install, and reload/restart the
# service so an upgrade always picks up the new binary (cargo-deb adds no systemd
# hooks by default).
set -e

if [ "$1" = "configure" ]; then
    # Install the default config only on first install; never overwrite an
    # existing one on upgrade (the config is not a dpkg conffile).
    if [ ! -e /etc/easylog/easylog.toml ] && [ -e /usr/share/easylog/easylog.toml ]; then
        mkdir -p /etc/easylog
        cp /usr/share/easylog/easylog.toml /etc/easylog/easylog.toml
    fi

    if [ -d /run/systemd/system ]; then
        systemctl daemon-reload >/dev/null 2>&1 || true
        # On upgrade ($2 = previous version) always restart so the new binary
        # takes effect — `restart` (not `try-restart`) also brings the service
        # back if it had crashed/stopped. Fresh installs leave enabling/starting
        # to the admin (`systemctl enable --now easylog`).
        if [ -n "$2" ]; then
            systemctl restart easylog.service >/dev/null 2>&1 || true
        fi
    fi
fi

exit 0
