# packaging/easylog.service
# systemd service unit for EasyLog (installed by the deb/rpm packages).
# Enable + start with:  sudo systemctl enable --now easylog
[Unit]
Description=EasyLog multi-log analyzer
Documentation=https://github.com/yarivha/EasyLog
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# EasyLog is a self-contained binary (templates + static assets are embedded),
# so no special working directory is needed.
Environment=EASYLOG_CONFIG=/etc/easylog/easylog.toml
ExecStart=/usr/bin/easylog
Restart=on-failure
RestartSec=2

# Runs as root: a syslog collector binds the privileged port 514 and this keeps
# the unit portable across hosts and containers (LXC, where DynamicUser and
# mount-namespace sandboxing are unreliable). StateDirectory creates and owns
# the database directory /var/lib/easylog. Lightweight, namespace-free hardening
# only, so the service starts in restricted container environments.
StateDirectory=easylog
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
